# rootfs/etc/uci-defaults/20-opkg-disable-https
#!/bin/sh

# Disable opkg https as no crypto libustream is installed
sed -i 's/https:/http:/g' /etc/opkg/distfeeds.conf